# airootfs/etc/skel/.config/autostart/prevail-lutris-optimizations.desktop
[Desktop Entry]
Name=Prevail-Lutris-Optimizations

Exec=prevail-lutris-optimizations

Terminal=false
Type=Application
Hidden=True