# airootfs/etc/skel/.config/autostart/set-avatar.desktop
[Desktop Entry]
Name=Live-Avatar

Exec=prevail-avatar

Terminal=false
Type=Application
Hidden=True