# airootfs/root/.config/autostart/set-avatar.desktop
[Desktop Entry]
Name=Wait for network to become available

Exec=prevail-avatar

Terminal=false
Type=Application
Hidden=True