# airootfs/etc/systemd/system/paccache.service
[Unit]
Description=Clean up old pacman packages

[Service]
Type=oneshot
ExecStart=/usr/bin/paccache -r
