# airootfs/etc/systemd/system/paccache.timer
[Unit]
Description=Run paccache daily

[Timer]
OnCalendar=daily
Persistent=true

[Install]
WantedBy=timers.target
