# airootfs/etc/xdg/autostart/network.desktop
[Desktop Entry]
Name=Network Check

Exec=network-check

Terminal=false
Type=Application
Hidden=True
