# airootfs/etc/xdg/autostart/user_services.desktop
[Desktop Entry]
Name=User Services

Exec=prevail-user-services

Terminal=false
Type=Application
Hidden=True
